# deploy/systemd/portfolio-weekly.service
[Unit]
Description=Portfolio Weekly Summary Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/portfolio-telegram-analytics
EnvironmentFile=/opt/portfolio-telegram-analytics/.env
ExecStart=/opt/portfolio-telegram-analytics/scripts/run_weekly.sh
NoNewPrivileges=true